# pkg/common/venv-salt-minion-postinstall.timer
[Unit]
Description=Clean old venv-salt-minion environment in 60 seconds

[Timer]
OnActiveSec=60

[Install]
WantedBy=timers.target

